# rclick-agent/dev.zwk.rclick-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>dev.zwk.rclick-agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/zwk/Applications/RClick Agent.app/Contents/MacOS/RClick Agent</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<!-- Aqua only: this is a GUI accessory app and must not be loaded at the
	     login window (pre-WindowServer), where it would fail to start and
	     crash-loop under KeepAlive. -->
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>ProcessType</key>
	<string>Interactive</string>
</dict>
</plist>
